# amazon-linux-2/nat/dbs-port-forwarding.service
[Unit]
Description=DBs Port Forwarding
[Service]
ExecStart=/home/ec2-user/dbs-port-forwarding
